# deploy/quadlet/imagecli-gen.container
[Unit]
Description=imageCLI generation NATS worker
Wants=lyra-nats.service
After=lyra-nats.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Container]
Image=ghcr.io/roxabi/imagecli-gen:staging
Label=io.containers.autoupdate=registry
ContainerName=imagecli-gen
Network=roxabi.network
# HuggingFace cache — bind-mount host cache so FLUX/SD3.5 weights are shared with
# other CLIs (llmCLI, etc.) and directly inspectable on the host (~10–14 GB).
# UserNS=keep-id:uid=1503,gid=1503 handles write permissions; no :z needed (no SELinux).
Volume=%h/.cache/huggingface:/home/imagecli/.cache/huggingface
# Explicit HF_HOME guards against XDG_CACHE_HOME quirks inside the container.
Environment=HF_HOME=/home/imagecli/.cache/huggingface
# Generated PNGs (nats_out/) — named volume keeps Syncthing-replicated outputs
# off the read-only container fs. Hub mounts the same volume read-only when
# delivering by file_path (output_mode=file path of the v0.4.x ImageResponse).
Volume=imagecli-out:/home/imagecli/.roxabi/imagecli
Secret=imagecli-nats,type=mount,target=/home/imagecli/.config/imagecli/nkeys/seed.seed,mode=0400,uid=1503,gid=1503
AddDevice=nvidia.com/gpu=all
# UID 1503 fixed in image; anchor mapping so the mounted nkey (uid=1503) is readable.
UserNS=keep-id:uid=1503,gid=1503
# Big-bang NATS consolidation — connect to shared lyra-nats on roxabi.network.
Environment=NATS_URL=nats://lyra-nats:4222
Environment=NATS_NKEY_SEED_PATH=/home/imagecli/.config/imagecli/nkeys/seed.seed
# Engine override without rebuilding the image (see Dockerfile.gen CMD).
Environment=IMAGECLI_ENGINE=flux2-klein
HealthCmd=pgrep -f "imagecli nats-serve"
HealthInterval=30s
NoNewPrivileges=true
DropCapability=all

[Service]
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
